# support/supervisor/my-app.conf
[program:my-app]
command=/opt/my-app/bin/my-app
